# debian/weatherstation-checksensor.cron.d
#
# Regular cron jobs for the weatherstation package
#
*/15 * 	* * * 	weatherstation	/usr/bin/checksensor
