# debian/weatherstation-chart.cron.d
#
# Regular cron jobs for the weatherstation package
#
0,15,30,45 * 	* * * 	weatherstation	/usr/bin/weather_chart
